# calorie-bot-env-watch.service
[Unit]
Description=Restart calorie-bot when .env changes

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart calorie-bot
